# rootfs/etc/services.d/cloudflared/run
#!/usr/bin/with-contenv sh
cloudflared tunnel --loglevel run ${CLOUDFLARE_TUNNEL_UUID}